# deploy/supervisor/linscio-api.conf
[program:linscio-api]
command=/www/wwwroot/linscio-portal/api/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
directory=/www/wwwroot/linscio-portal/api
user=www
autostart=true
autorestart=true
startretries=3
stdout_logfile=/www/wwwroot/linscio-release/logs/api.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/www/wwwroot/linscio-release/logs/api_error.log
environment=ENV="production"
