# deploy/supervisor/linscio-api.conf
[program:linscio-api]
command=/home/linscio/api/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
directory=/home/linscio/api
user=linscio
autostart=true
autorestart=true
startretries=3
stdout_logfile=/home/linscio/logs/api.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/home/linscio/logs/api_error.log
environment=ENV="production"
